# xgbapi/xgbapi.service
[Unit]
Description=XGBoost Predict API (FastAPI/Uvicorn)
After=network.target

[Service]
User=centos
Group=centos
fayenp-codex/generate-fastapi-xgboost-api-on-centos-stream-9
WorkingDirectory=%h/xgbapi
Environment=PYTHONPATH=%h/xgbapi
EnvironmentFile=%h/xgbapi/config/.env
ExecStart=%h/xgbapi/.venv/bin/uvicorn app.main:app --host ${API_HOST} --port ${API_PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
